# docker-compose/assets/build/etc/profile.d/workspace_dir.sh
export WORKSPACE_DIR=$HOME/workspace
export ASSETS_DIR=$WORKSPACE_DIR/assets/docker/workspace
export RUNTIME_DIR=$HOME/runtime